# usb_debug.service
[Unit]
Description=start usb debugging

[Service]
ExecStart=bash /usb_debug.sh

